# deploy/aiqa-server.service
[Unit]
Description=AIQA Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/aiqa/server
Environment="NODE_ENV=production"
EnvironmentFile=/opt/aiqa/server/.env
ExecStart=/usr/bin/env node dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aiqa-server

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

